# roles/nuclei_scanner/templates/disable-scan-timer@.service.j2
[Unit]
Description=Disable Nuclei timer if it fails

[Service]
Type=oneshot
ExecStart=/bin/systemctl disable nuclei-%i.timer --now

[Install]
WantedBy=default.target
