# roles/nuclei_scanner/templates/nuclei-scan.service.j2
[Unit]
Description=Nuclei scanning service - {{ type }}
OnFailure=nuclei-disable-scan-timer@{{ type }}.service

[Service]
Type=exec

WorkingDirectory=/home/ubuntu/vulnman-nuclei/

User=ubuntu
Group=ubuntu

ExecStartPre=git pull
ExecStartPre=poetry install
ExecStartPre=nuclei -update 
ExecStartPre=nuclei -update-templates
ExecStartPre=send-notification "Nuclei" "Start - {{ type }}"

ExecStart=poetry run nuclei-scan-runner --config ./configs/{{ type }}-config.toml --no-update-nuclei --no-update-templates
ExecStopPost=send-notification "Nuclei" "Stop - {{ type }}"
