# roles/domain_discovery/templates/domain_discovery.service.j2
[Unit]
Description=Timer for domain discovery

[Service]
Type=exec
ExecStart=domain_discovery_start