# roles/prod_database/templates/backup_db.service.j2
[Unit]
Description=Postgres database backup service

[Service]
Type=oneshot

ExecStart={{ prod_database_path }}/backup_db.sh

[Install]
WantedBy=multi-user.target
